# package/my/luci-theme-darkpro/root/etc/uci-defaults/30_luci-theme-darkpro
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Darkpro=/luci-static/darkpro
	set luci.main.mediaurlbase=/luci-static/darkpro
	commit luci
EOF
exit 0
